# shared/cwsd.service
[Unit]
Description=CW sender daemon
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cwsd
ExecReload=/bin/kill -HUP $MAINPID
RestartSec=60
Restart=always
User=benny
Group=users

# Allow the keyer thread to use SCHED_FIFO and lock its memory so element timing
# is not stretched by ordinary scheduler preemption / page faults.
AmbientCapabilities=CAP_SYS_NICE CAP_IPC_LOCK
LimitRTPRIO=20
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target

